# wikiwalk.path
[Unit]
Description=Watch for changes to wikiwalk server

[Path]
PathChanged=/usr/local/bin/wikiwalk
Unit=wikiwalk.service

[Install]
WantedBy=multi-user.target
